# obnam/templates/etc/cron.d/obnam-host.erb
# cron entry for obnam host profile
#
# We do regular backup and some housekeeping to keep the backup storage usage under control.

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
MAILTO=root
HOME=/

OBNAM="obnam --quiet --config /etc/obnam/host.profile"
<%= cron_minute %> <%= cron_hour %> * * * root $OBNAM backup && $OBNAM forget
0 13 * * * root $OBNAM force-lock
